# v2/manifests/packages/shared.list
# CeraLive — Canonical SHARED runtime package list
#
# Arch-independent Debian packages every CeraLive device installs, regardless of
# SoC/board. This is the SINGLE source of truth, reconciled from three diverging
# legacy sources (Task 18):
#   1. configs/base/ceraui-base.conf        (BASE_PACKAGES + STREAMING_PACKAGES)
#   2. userpatches/customize-image.sh        (the set actually installed in v1)
#   3. configs/devices/{rock5bplus,orangepi5plus}.conf (BOARD_PACKAGES — BSP only)
#
# NOT here (by design, see removed.md for the full accounting):
#   - BSP / HW-accel (kernel, DTB, U-Boot, firmware, gstreamer1.0-rockchip1,
#     rockchip-multimedia-config) → v2/manifests/families/rk3588.yaml
#   - RK3588-specific runtime deltas → rk3588.delta.list (currently empty)
#   - First-party .debs (the CERALIVE/srt runtime, cerastream, srtla-send-rs, and
#     CeraUI) → installed in the APPLICATION layer (Stage 3), never fetched from
#     Debian by this manifest. The SRT package replaces Debian's TLS-flavor
#     variants in the same offline local transaction; see LAYER-MAP.md §libsrt.
#   - Development tooling (build-essential, git, vim, …) → development variant only
#   - openssh-server, systemd, udev, dbus → already in the mkosi base layer
#
# Format: one package per line; `#` comments justify non-obvious entries.
# The mkosi runtime layer (v2/mkosi/mkosi.images/runtime/mkosi.conf) installs
# exactly this set; lib/parity-check.sh diffs the built rootfs against it.

# --- TLS / apt plumbing -------------------------------------------------------
ca-certificates          # TLS trust store; required for https apt + apt.ceralive.tv mTLS
apt-transport-https      # https:// apt transport for the first-party repo
openssl                  # CLI for the first-boot self-signed TLS cert (ceralive-tls-firstboot).
                         # Priority: important on Debian (usually already present), pinned
                         # explicitly so cert generation never depends on a transitive pull.

# --- CeraUI TLS front (Task 15, SC3) ------------------------------------------
nginx-light              # terminates HTTPS on 443 and reverse-proxies to the CeraUI
                         # backend on 127.0.0.1:80 (WebSocket-upgrade aware, EC6). Port 80
                         # is LEFT to the backend (no redirect). nginx-light is the smallest
                         # viable nginx flavour — it carries the proxy + ssl modules without
                         # the mail/stream/extras bulk of nginx-full. Cert is per-device
                         # self-signed on /data; see v2/docs/size-notes.md for the size delta.

# --- Network management -------------------------------------------------------
network-manager          # primary connection manager (WiFi/Ethernet/modem); CeraUI drives it
systemd-resolved         # DNS resolver; NetworkManager uses dns=systemd-resolved (postinst-lib.sh). On Debian bookworm this is a SEPARATE package (was bundled in systemd on v1's Armbian base).
dnsmasq                  # DHCP/DNS for hostapd AP mode + local resolution
iproute2                 # `ip`/rt_tables — SRTLA source-policy routing (customize-image.sh:294)
ipcalc                   # subnet math in the NM dispatcher (customize-image.sh:358)

# --- USB modem bring-up -------------------------------------------------------
usb-modeswitch           # flips USB modems out of storage mode into modem mode —
                         # required for USB LTE/5G modem enumeration (core datapath).
                         # NOTE: audio (alsa-utils/pulseaudio) and the rest of the USB/HW
                         # introspection set (usbutils/pciutils/lsof/i2c-tools/can-utils)
                         # moved to the debug add-on (T17) — see removed.md.

# --- GStreamer core (streaming pipeline; HW-accel plugin is in the family manifest)
gstreamer1.0-tools       # gst-launch/gst-inspect — pipeline plumbing
gstreamer1.0-plugins-base
gstreamer1.0-plugins-good
gstreamer1.0-plugins-bad
gstreamer1.0-plugins-ugly
gstreamer1.0-libav       # libav-backed encode/decode fallback
v4l-utils                # v4l2-ctl + media-ctl for HDMI/UVC capture device control
ffmpeg                   # transcode/probe utility for the pipeline

# --- Modems (LTE/5G on every board) -------------------------------------------
modemmanager             # cellular modem management daemon
libqmi-utils             # QMI modem control (Qualcomm)
libmbim-utils            # MBIM modem control

# --- WiFi ---------------------------------------------------------------------
wireless-tools           # legacy iw* tooling
wpasupplicant            # WPA client
wireless-regdb           # wireless regulatory DB (/lib/firmware/regulatory.db + .p7s) — the
                         # kernel cfg80211 subsystem loads it at boot to establish a usable
                         # regulatory domain. It is wpasupplicant's `Recommends:`, so it is
                         # NOT pulled transitively (this list installs with
                         # --no-install-recommends); named explicitly or cfg80211 boots with
                         # "Direct firmware load for regulatory.db failed with error -2" and
                         # NetworkManager finds no usable WiFi interface (real-HW UART, 2026-07-16).
hostapd                  # AP mode (device hotspot for config)

# --- mDNS / discovery ---------------------------------------------------------
avahi-daemon             # .local hostname advertisement (find the box on the LAN)
libnss-mdns              # NSS plugin so the box can resolve other .local hosts

# --- Time ---------------------------------------------------------------------
chrony                   # NTP client; enabled service in customize-image.sh:521.
                         # NOTE: v1 enabled the unit but never apt-installed it — fixed here.

# --- CPU / power --------------------------------------------------------------
cpufrequtils             # CPU governor control for sustained encode performance

# --- Net glue (runtime-critical) ----------------------------------------------
rsync                    # file sync for dev-push/dev-sync live-reload loop (v2/docs/dev-loop.md)
curl                     # http client — CeraUI/update flows + on-device `curl localhost`
                         # health probe (backend observability.ts). Runtime-critical: KEPT.
wget                     # http downloader (update/fetch infra)
                         # NOTE: net DIAGNOSTICS (socat/netcat-openbsd/iperf3) moved to the
                         # debug add-on (T17) — verified unused by CeraUI/streaming runtime;
                         # the chaos harness (tools/chaos) ships its own. See removed.md.
sudo                     # required by CeraUI's narrow /etc/sudoers.d/ceralive-addon-helper
                         # drop-in; without the sudo binary add-on enable/disable is dead code.

# --- Console (HDMI readability at boot) ---------------------------------------
kbd                      # setfont — set console font at boot for HDMI readability at any resolution
fonts-terminus           # large bitmap fonts (ter-v32n etc.) readable at 4K on fbcon
                         # NOTE: operational monitoring (htop/iotop/nethogs/vnstat) and the
                         # nano editor moved to the debug add-on (T17) — see removed.md.

# --- OS-update infrastructure (Stage 2; RUNTIME layer, arch-independent) -----
# The device SRT runtime is CeraLive's `libsrt1.5-ceralive` first-party package,
# staged and installed by the App layer. Do not add a Debian `libsrt1.5-*` flavor
# here: GStreamer and cerastream must resolve one forked shared object. The LAN SRT
# ingest listener is served by MediaMTX, so no srt-tools package is needed.
nftables                 # ships /usr/sbin/nft — loads the LAN-ingest ingress firewall
                         # (ceralive-ingest-firewall.service): DROP inbound :1935 (RTMP) +
                         # :4001 (SRT) on modem/WWAN/ppp uplink classes so the
                         # UNAUTHENTICATED v1 MediaMTX ingest gateway is reachable from
                         # LAN/hotspot ONLY, never from a WAN uplink (DEFERRED.md items 7 & 8).
                         # libnftables1 is already pulled by the base (NetworkManager dep), so
                         # the net add is just the nft binary + libmnl/libnftnl (~1.5 MB — see
                         # size-notes.md §7).
rauc                     # A/B update client (bookworm `rauc` 1.8-2). decisions.md T5.
rauc-service             # ships /lib/systemd/system/rauc.service (D-Bus daemon). bookworm
                         # SPLIT the daemon out of `rauc` (binary-only); the runtime
                         # postinst enables rauc.service, so the service pkg is required.
squashfs-tools           # ships /usr/bin/unsquashfs — rauc shells out to it to extract the
                         # manifest (and rootfs image) from a plain-format .raucb at
                         # `rauc info`/`rauc install` time. Absent it, install fails after
                         # signature verification with "Failed to start unsquashfs: ... No
                         # such file or directory" (real Rock 5B+ hardware). build-time
                         # mksquashfs runs on the HOST/CI, so this runtime gap was invisible.
e2fsprogs                # ships /sbin/mkfs.ext4 — rauc shells out to it to format the target
                         # ext4 slot before copying the new rootfs image during the slot-write
                         # phase of `rauc install`, AFTER signature + manifest checks pass.
                         # Without it, real Rock 5B+ hardware reported exactly: "LastError: Installation error: Failed updating slot rootfs.1: failed to start mkfs.ext4: Failed to execute child process 'mkfs.ext4' (No such file or directory)".
                         # build-time tooling never needed mkfs.ext4, so this runtime-only gap was
                         # invisible; installing e2fsprogs made the REAL end-to-end OTA install
                         # complete successfully (signature verified, slot B activated and rebooted
                         # healthy), proving the fix on physical hardware.
u-boot-tools             # fw_setenv/fw_getenv + mkimage — RAUC bootloader env / boot
                         # selector. Generic Debian pkg (present on amd64 too); the
                         # board-specific U-Boot BLOB is a Platform-layer BSP concern.
# rauc-hawkbit-updater   # Stage 7 OTA (task 41). DELIBERATELY a comment, NOT an
                         # active apt line: 1.4 is NOT in bookworm main (only sid/forky
                         # 1.4-1) → it is a bookworm BACKPORT .deb (1.4-1~bpo12+1),
                         # pre-staged like the first-party .debs and installed by the
                         # RUNTIME postinst (setup_hawkbit_updater → apt-get install
                         # ./<file>.deb, which pulls libcurl4/libglib2.0-0/libjson-glib
                         # from bookworm main). An active line here would break the
                         # whole-list `apt-get install` on any host without the backport.
                         # Config/enrollment: mkosi/runtime/hawkbit-updater/. decisions.md Task 5.
