# v2/manifests/packages/shared.list
# CeraLive — Canonical SHARED runtime package list
#
# Arch-independent Debian packages every CeraLive device installs, regardless of
# SoC/board. This is the SINGLE source of truth, reconciled from three diverging
# legacy sources (Task 18):
#   1. configs/base/ceraui-base.conf        (BASE_PACKAGES + STREAMING_PACKAGES)
#   2. userpatches/customize-image.sh        (the set actually installed in v1)
#   3. configs/devices/{rock5bplus,orangepi5plus}.conf (BOARD_PACKAGES — BSP only)
#
# NOT here (by design, see removed.md for the full accounting):
#   - BSP / HW-accel (kernel, DTB, U-Boot, firmware, gstreamer1.0-rockchip1,
#     rockchip-multimedia-config) → v2/manifests/families/rk3588.yaml
#   - RK3588-specific runtime deltas → rk3588.delta.list (currently empty)
#   - First-party .debs (cerastream, srtla, CeraUI, and the CERALIVE/srt fork) →
#     installed in the APPLICATION layer (Stage 3), never apt packages. NOTE the
#     SYSTEM libsrt transport library (Debian `libsrt1.5-openssl`) IS here, in the
#     OS-update infra section below — it rides the RAUC OS slot so app sysexts stay
#     small and a libsrt bump is atomic. The CERALIVE/srt .deb (first-party fork)
#     is separate and lands in the App layer (see LAYER-MAP.md §libsrt).
#   - Development tooling (build-essential, git, vim, …) → development variant only
#   - openssh-server, systemd, udev, dbus → already in the mkosi base layer
#
# Format: one package per line; `#` comments justify non-obvious entries.
# The mkosi runtime layer (v2/mkosi/mkosi.images/runtime/mkosi.conf) installs
# exactly this set; lib/parity-check.sh diffs the built rootfs against it.

# --- TLS / apt plumbing -------------------------------------------------------
ca-certificates          # TLS trust store; required for https apt + apt.ceralive.tv mTLS
apt-transport-https      # https:// apt transport for the first-party repo

# --- Network management -------------------------------------------------------
network-manager          # primary connection manager (WiFi/Ethernet/modem); CeraUI drives it
systemd-resolved         # DNS resolver; NetworkManager uses dns=systemd-resolved (postinst-lib.sh). On Debian bookworm this is a SEPARATE package (was bundled in systemd on v1's Armbian base).
dnsmasq                  # DHCP/DNS for hostapd AP mode + local resolution
iproute2                 # `ip`/rt_tables — SRTLA source-policy routing (customize-image.sh:294)
ipcalc                   # subnet math in the NM dispatcher (customize-image.sh:358)

# --- USB modem bring-up -------------------------------------------------------
usb-modeswitch           # flips USB modems out of storage mode into modem mode —
                         # required for USB LTE/5G modem enumeration (core datapath).
                         # NOTE: audio (alsa-utils/pulseaudio) and the rest of the USB/HW
                         # introspection set (usbutils/pciutils/lsof/i2c-tools/can-utils)
                         # moved to the debug add-on (T17) — see removed.md.

# --- GStreamer core (streaming pipeline; HW-accel plugin is in the family manifest)
gstreamer1.0-tools       # gst-launch/gst-inspect — pipeline plumbing
gstreamer1.0-plugins-base
gstreamer1.0-plugins-good
gstreamer1.0-plugins-bad
gstreamer1.0-plugins-ugly
gstreamer1.0-libav       # libav-backed encode/decode fallback
v4l-utils                # v4l2-ctl + media-ctl for HDMI/UVC capture device control
ffmpeg                   # transcode/probe utility for the pipeline

# --- Modems (LTE/5G on every board) -------------------------------------------
modemmanager             # cellular modem management daemon
libqmi-utils             # QMI modem control (Qualcomm)
libmbim-utils            # MBIM modem control

# --- WiFi ---------------------------------------------------------------------
wireless-tools           # legacy iw* tooling
wpasupplicant            # WPA client
hostapd                  # AP mode (device hotspot for config)

# --- mDNS / discovery ---------------------------------------------------------
avahi-daemon             # .local hostname advertisement (find the box on the LAN)
libnss-mdns              # NSS plugin so the box can resolve other .local hosts

# --- Time ---------------------------------------------------------------------
chrony                   # NTP client; enabled service in customize-image.sh:521.
                         # NOTE: v1 enabled the unit but never apt-installed it — fixed here.

# --- CPU / power --------------------------------------------------------------
cpufrequtils             # CPU governor control for sustained encode performance

# --- Net glue (runtime-critical) ----------------------------------------------
rsync                    # file sync for dev-push/dev-sync live-reload loop (v2/docs/dev-loop.md)
curl                     # http client — CeraUI/update flows + on-device `curl localhost`
                         # health probe (backend observability.ts). Runtime-critical: KEPT.
wget                     # http downloader (update/fetch infra)
                         # NOTE: net DIAGNOSTICS (socat/netcat-openbsd/iperf3) moved to the
                         # debug add-on (T17) — verified unused by CeraUI/streaming runtime;
                         # the chaos harness (tools/chaos) ships its own. See removed.md.

# --- Console (HDMI readability at boot) ---------------------------------------
kbd                      # setfont — set console font at boot for HDMI readability at any resolution
fonts-terminus           # large bitmap fonts (ter-v32n etc.) readable at 4K on fbcon
                         # NOTE: operational monitoring (htop/iotop/nethogs/vnstat) and the
                         # nano editor moved to the debug add-on (T17) — see removed.md.

# --- OS-update + SRT transport infra (Stage 2; RUNTIME layer, arch-independent) -
# Added in Task 19 (layering). These belong in the RUNTIME OS slot — NOT the App
# sysext layer — so they update atomically via RAUC and are identical on every
# architecture (rk3588, x86). See LAYER-MAP.md and decisions.md (Task 5 RAUC pins).
libsrt1.5-openssl        # SYSTEM SRT transport library that the streaming stack (srtla) links
                         # at runtime. In the OS slot (not app sysext) → app images
                         # stay small + a libsrt bump is atomic via the RAUC slot.
                         # bookworm provides 1.5.x; the OpenSSL backend matches the
                         # CERALIVE build. (The first-party CERALIVE/srt FORK .deb is
                         # a SEPARATE App-layer artifact — see LAYER-MAP.md §libsrt.)
rauc                     # A/B update client (bookworm `rauc` 1.8-2). decisions.md T5.
rauc-service             # ships /lib/systemd/system/rauc.service (D-Bus daemon). bookworm
                         # SPLIT the daemon out of `rauc` (binary-only); the runtime
                         # postinst enables rauc.service, so the service pkg is required.
u-boot-tools             # fw_setenv/fw_getenv + mkimage — RAUC bootloader env / boot
                         # selector. Generic Debian pkg (present on amd64 too); the
                         # board-specific U-Boot BLOB is a Platform-layer BSP concern.
# rauc-hawkbit-updater   # Stage 7 OTA (task 41). DELIBERATELY a comment, NOT an
                         # active apt line: 1.4 is NOT in bookworm main (only sid/forky
                         # 1.4-1) → it is a bookworm BACKPORT .deb (1.4-1~bpo12+1),
                         # pre-staged like the first-party .debs and installed by the
                         # RUNTIME postinst (setup_hawkbit_updater → apt-get install
                         # ./<file>.deb, which pulls libcurl4/libglib2.0-0/libjson-glib
                         # from bookworm main). An active line here would break the
                         # whole-list `apt-get install` on any host without the backport.
                         # Config/enrollment: mkosi/runtime/hawkbit-updater/. decisions.md Task 5.
