# v2/mkosi/runtime/ceralive-ssh-firstboot.service
[Unit]
Description=CeraLive SSH hardening, persistent identity, and CI-key guard
Documentation=file:///usr/share/doc/ceralive/ssh-hardening.md
# Task 10 / SC4. Four scoped actions (see the script header): regenerate the
# baked shared SSH host keys into a per-device identity, drop in
# PermitRootLogin prohibit-password, and expire the ceralive password once so the
# first interactive login forces a change, and link the empty /data-backed
# authorized-key store. Nothing more (no fail2ban/UFW/auditd).
#
# MUST complete before sshd accepts any connection on every boot. Initial boot
# replaces the shared baked host-key fingerprint; later boots purge unauthorized
# CI keys before either the classic service or socket activation can accept a
# connection.
#
# ssh.socket is ordered Before=sockets.target (early boot, before basic.target),
# so a guard ordered Before=ssh.socket must NOT inherit the default
# After=basic.target — that forms an unbreakable ssh.socket<->sockets.target
# ordering cycle (systemd deletes ssh.socket's start job and SSH never comes up).
# So opt out of default deps, but re-add the SAFE half by hand. This script does
# sysinit-phase-dependent work (ssh-keygen -A needs the seeded RNG; the
# authorized-key store setup resolves the ceralive user/group via NSS and chowns
# to it; sshd -t needs tmpfiles/udev-provided runtime paths), so it MUST run after
# sysinit.target. proof-11 proved that dropping After=sysinit.target (a side effect
# of DefaultDependencies=no) let this unit race ahead of systemd-tmpfiles/sysusers/
# udev and fail under `set -euo pipefail`, taking ssh.service/ssh.socket down with
# "Dependency failed". After=sysinit.target restores every sysinit-phase
# prerequisite WITHOUT re-adding After=basic.target — sysinit.target is ordered
# before sockets.target, so it can never close the ssh.socket loop.
DefaultDependencies=no
Requires=sysinit.target
After=sysinit.target
Conflicts=shutdown.target
Before=shutdown.target
After=local-fs.target data.mount
Before=ssh.service ssh.socket
# Persisted per-device host identity lives on /data (survives reboots + A/B OTA
# slot swaps). data.mount is a soft ordering: images with no /data partition fall
# back to /etc/ceralive inside the script (mountpoint-guarded), so do NOT make the
# mount a hard requirement here.
ConditionPathExists=/etc/ssh

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ceralive-ssh-firstboot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
RequiredBy=ssh.service ssh.socket
