# v2/mkosi/runtime/ceralive-ssh-firstboot.service
[Unit]
Description=CeraLive first-boot SSH hardening — per-device host keys, root password login disabled, forced password change
Documentation=file:///usr/share/doc/ceralive/ssh-hardening.md
# Task 10 / SC4. Three scoped actions (see the script header): regenerate the
# baked shared SSH host keys into a per-device identity, drop in
# PermitRootLogin prohibit-password, and expire the ceralive password once so the
# first interactive login forces a change. Nothing more (no fail2ban/UFW/auditd).
#
# MUST run before sshd accepts ANY connection on first boot, otherwise the very
# first client would pin the shared baked host-key fingerprint (a MITM hazard)
# and root could still authenticate by password. Ordered Before both the classic
# service and the socket so it wins regardless of which sshd activation the image
# uses.
After=local-fs.target data.mount
Before=ssh.service ssh.socket
# Persisted per-device host identity lives on /data (survives reboots + A/B OTA
# slot swaps). data.mount is a soft ordering: images with no /data partition fall
# back to /etc/ceralive inside the script (mountpoint-guarded), so do NOT make the
# mount a hard requirement here.
ConditionPathExists=/etc/ssh

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ceralive-ssh-firstboot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
