# v2/mkosi/runtime/ingest-firewall/ceralive-ingest-firewall.service
[Unit]
Description=CeraLive LAN-ingest ingress firewall (drop :1935/:4001 on WAN/modem uplinks)
Documentation=https://www.netfilter.org/projects/nftables/
# INGRESS BOUNDARY (security, v1). Installs the nftables ruleset that keeps the
# UNAUTHENTICATED LAN ingest gateway — ceralive-rtmp-gateway.service (MediaMTX,
# RTMP :1935 + SRT :4001) — reachable ONLY from LAN/hotspot, never from a
# modem/WWAN/ppp WAN uplink. In v1 the gateway carries no publish password or SRT
# passphrase (LAN-scoped trust posture — see v2/docs/DEFERRED.md items 7 & 8); this
# firewall is what enforces that boundary. Order this BEFORE the gateway so the
# boundary is in place before either listener opens its socket.
DefaultDependencies=no
After=nftables.service systemd-sysctl.service
Before=ceralive-rtmp-gateway.service network-pre.target shutdown.target
Wants=network-pre.target
Conflicts=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
# iifname wildcard matching binds by NAME at packet time, so the rules are correct
# even for modem/wifi interfaces that appear after boot — no network-online wait.
ExecStart=/usr/sbin/nft -f /etc/ceralive/ingest-firewall.nft
ExecReload=/usr/sbin/nft -f /etc/ceralive/ingest-firewall.nft
# Tear down only OUR table on stop (leading '-' = ignore if already absent); never
# flush the whole ruleset, so any other firewall table is left intact.
ExecStop=-/usr/sbin/nft delete table inet ceralive_ingest_fw

[Install]
WantedBy=multi-user.target
